# wireguard-qr-manager.service
[Unit]
Description=WireGuard QR Manager - FastAPI Service
After=network.target wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/wireguard-qr-manager
EnvironmentFile=/opt/wireguard-qr-manager/.env

# Python virtual environment
ExecStart=/opt/wireguard-qr-manager/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 6000

# Restart policy
Restart=always
RestartSec=5

# Security (adjust as needed - root required for wg commands)
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-qr-manager

[Install]
WantedBy=multi-user.target
